# systemd/email-tunnel@.service
[Unit]
Description=Hermes Email Tunnel for %i
After=network-online.target mihomo.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/.local/bin/hermes-email-tunnel %i
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
